# templates/devpi.conf
[program:devpi]
command = devpi-server --serverdir /srv/devpi --host 0.0.0.0
user = vagrant
autostart = True
